# eg25-autorestart/eg25-autorestart.service
[Unit]
Description=Quectel EG25 autorestart
RefuseManualStop=true

[Service]
Type=oneshot
ExecStart=/usr/lib/eg25-autorestart/eg25-autorestart.sh
