# eg25-autorestart/eg25-autorestart.timer
[Unit]
Description=Quectel EG25 autorestart timer

[Timer]
OnStartupSec=5min
OnUnitInactiveSec=5min

[Install]
WantedBy=timers.target
